# ziz-main.service
[Unit]
Description=Run zizzania daemonized

[Service]
Restart=always

EnvironmentFile=/opt/zizzania/environment
ExecStart=/usr/local/bin/zizzania -i $ADAPTER -M -d 10 -w /tmp/zizzania.pcap
#StandardOutput=file:/tmp/zizzania.stdout

[Install]
WantedBy=multi-user.target
